# files/etc/etc/crontabs/root
* */8 * * * /sbin/free.sh >/dev/null 2>&1
